# debian/lbnamed.init
#! /bin/sh
### BEGIN INIT INFO
# Provides:          lbnamed
# Required-Start:    $syslog $remote_fs
# Required-Stop:     $syslog $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Balancing DNS written in Perl
# Description:       Balancig DNS launcher
### END INIT INFO

cd /var/lib/lbnamed
/sbin/lbnamed -l /var/log/lbnamed.log
